# resources/docker/etc/supervisor.d/crond.conf
[program:crond]
command = /usr/sbin/crond -f
autostart = true
autorestart = true
